# packages/arch/shrt.service
[Unit]
Description=Shrt service
After=network.target

[Service]
User=shrt
Group=shrt
ExecStart=/usr/bin/shrt
Restart=on-failure

[Install]
WantedBy=multi-user.target
